# ansible/roles/hadoop/templates/hadoop-hdfs.service.j2
[Service]
Type=simple
User=vagrant
Group=vagrant
ExecStart=/opt/hadoop/bin/hdfs namenode
Restart=always
RestartSec=10